# newroot/etc/sv/getty-6/run
#!/bin/sh

exec 2>&1

echo "getty-6!!.."

setsid /bin/sh -c 'exec login <> /dev/tty6 >&0 2>&1'

echo "..out of GETTY-6!!!!!!!!!!!!!11111111"
