# newroot/etc/sv/tester/run
#!/bin/sh

echo "ohhhh!!.."

exec /tester.sh 
