# newroot/etc/sv/tester/run
#!/bin/sh

exec 2>&1

exec echo "ohhhh!!.."

#exec /tester.sh 
exec setsid /sbin/agetty --nohostname tty5 38400 linux
