# newroot/etc/sv/tester/run
#!/bin/sh

exec 2>&1

exec echo "in!!.."

#exec /tester.sh 
#exec setsid /sbin/agetty --nohostname tty5 38400 linux
exec setsid /sbin/getty tty6 38400 raaNux

exec echo "..out"
